# templates/supervisor/project.conf.erb
[program:project]
command=/home/obs/virtualenv/bin/uwsgi
  --virtualenv /home/obs/virtualenv
  --module app.wsgi:application
  --socket /var/uwsgi/project.socket
  --chmod-socket
  --master
  --processes 2
  --harakiri 300
  --max-requests 100
  --stats /var/uwsgi/project_stats.socket
directory=/home/obs/project
environment=DJANGO_SETTINGS_MODULE='app.settings'
user=obs
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/project.log
redirect_stderr=true
stopsignal=QUIT
